# kafka/config/supervisor_rtlogtrans_client.ini

[program:supervisor_rtlogtrans_client]
directory = /home/mtime/optools/rtlogtrans-service
command = /home/mtime/optools/rtlogtrans-service/env/bin/python rtlogtrans_client.py
autostart = true
autorestart = true
startsecs = 3
startretries = 3
stopwaitsecs = 5
user = root
redirect_stderr = true
stdout_logfile = /home/mtime/logs/supervisor/supervisor_rtlogtrans_client.log
stdout_logfile_maxbytes = 500MB
stdout_logfile_backups = 0
