# kafka/config/supervisor_rtlogtrans-service.ini

[program:supervisor_rtlogtrans-service]
directory = /home/mtime/optools/rtlogtrans-service
command = /home/mtime/optools/rtlogtrans-service/env/bin/python rtlogtrans.py
autostart = true
autorestart = true
startsecs = 3
startretries = 3
stopwaitsecs = 5
user = root
redirect_stderr = true
stdout_logfile = /home/mtime/logs/supervisor/supervisor_rtlogtrans-service.log
stdout_logfile_maxbytes = 500MB
stdout_logfile_backups = 0
